# roles/supervisor/templates/service.conf
[program:{{ service }}]
command={{ command }}            ; the program (relative uses PATH, can take args)
;priority=999                ; the relative start priority (default 999)
autostart=true              ; start at supervisord start (default: true)
autorestart=true            ; retstart at unexpected quit (default: true)
startsecs=5                ; number of secs prog must stay running (def. 10)
startretries=3              ; max # of serial start failures (default 3)
exitcodes=0,2               ; 'expected' exit codes for process (default 0,2)
stopsignal=TERM             ; signal used to kill process (default TERM)
stopwaitsecs=10             ; max num secs to wait before SIGKILL (default 10)
log_stdout=true             ; if true, log program stdout (default true)
log_stderr=true             ; if true, log program stderr (def false)
stdout_logfile=/var/log/supervisor/{{ service }}.log
stderr_logfile=/var/log/supervisor/{{ service }}.log
{% if env is defined %}
environment={{ env }}
{% endif %}
